# root/etc/uci-defaults/30_luci-theme-glassnova
#!/bin/sh
uci -q batch <<'UCI'
set luci.themes.GlassNova='/luci-static/glassnova'
commit luci
UCI
exit 0
